# cockpit-api.service
[Unit]
Description=Cockpit API — Issue-Driven Spec Pipeline
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=%i
WorkingDirectory=/home/%i/repos/cockpit/backend

# Extend PATH so the service can find claude (~/.local/bin) and other user tools
Environment=PATH=/home/%i/.local/bin:/usr/local/bin:/usr/bin:/bin

# Load env vars from .env file in the repo root
EnvironmentFile=/home/%i/repos/cockpit/.env

# Use the venv created by: cd backend && python -m venv .venv && .venv/bin/pip install -r requirements.txt
ExecStart=/home/%i/repos/cockpit/backend/.venv/bin/uvicorn main:app --host 0.0.0.0 --port 8001

Restart=on-failure
RestartSec=5

# Give uvicorn time to shut down cleanly
TimeoutStopSec=30

StandardOutput=journal
StandardError=journal
SyslogIdentifier=cockpit-api

[Install]
WantedBy=multi-user.target
